# systemctl-services/owfs.service
[Unit]
Description=1-wire filesystem FUSE mount
Documentation=man:owfs(1)

[Service]
Type=notify
NotifyAccess=all
ExecStart=/usr/bin/owfs -c /etc/owfs.conf -uall --allow_other %t/owfs
ExecStop=/usr/bin/umount %t/owfs
RuntimeDirectory=owfs

[Install]
WantedBy=multi-user.target
